# profile.d/x11vnc-user.sh
x11vnc -auth guess -forever -loop -noxdamage -repeat -rfbauth /opt/sicoob/passwdvnc.pwd -rfbport 5900 -shared -o /var/log/x11vnc.log
