# profile.d/x11vnc-user.sh
/usr/bin/x11vnc -forever -noxdamage -repeat -rfbauth /opt/sicoob/passwdvnc.pwd -autoport 5900 -shared -o .local/share/x11vnc.log -q &
